# infra/pi/mc-orchestrator.service
[Unit]
Description=Minecraft LLM Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/mc
ExecStart=/usr/bin/env npm run --workspace @mc/orchestrator start -- --config /etc/mc-orchestrator/config.yaml
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGTERM
StandardOutput=append:/var/log/mc-orchestrator/service.log
StandardError=append:/var/log/mc-orchestrator/service.err

[Install]
WantedBy=multi-user.target
